# toast/miniupnpd/miniupnpd.initd
#!/sbin/openrc-run
command="/usr/sbin/miniupnpd"
pidfile="/run/${RC_SVCNAME}.pid"
command_args="$ARGS"

depend() {
	need net firewall
	use minissdpd
}

start_pre() {
    /etc/miniupnpd/iptables_init.sh
}

stop_post() {
    /etc/miniupnpd/iptables_removall.sh
}
