# toast/miniupnpd/miniupnpd.initd
#!/sbin/openrc-run
command="/usr/sbin/miniupnpd"
pidfile="/run/${RC_SVCNAME}.pid"
command_args="$ARGS -P $pidfile"

depend() {
    need firewall
    use minissdpd
}

start_pre() {
    /etc/miniupnpd/iptables_init.sh
}

stop_post() {
    /etc/miniupnpd/iptables_removeall.sh
}
